# services/headphones/run
#!/bin/bash

exec /sbin/setuser abc python /opt/headphones/Headphones.py -p 8181 --datadir=/config --nolaunch
