# root/etc/s6-overlay/s6-rc.d/init-headphones-config/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

# copy config
if [[ ! -f /config/config.ini ]]; then
    cp /defaults/config.ini /config/config.ini
fi

# permissions
lsiown -R abc:abc \
    /config
