# root/etc/services.d/headphones/run
#!/usr/bin/with-contenv bash

exec \
	s6-setuidgid abc python /app/headphones/Headphones.py \
	--datadir=/config
