# services/tail/run
#!/bin/bash

exec /sbin/setuser abc tail -F /config/logs/headphones.log